# .docker/sync.crontab
0 * * * * cd /home/jurien/varrock-stray-dog && /snap/bin/docker-compose up sync >>/tmp/sync.log 2>&1